# deploy/evomesh.service
# This file is a TEMPLATE. setup.sh generates the actual service file
# with correct paths and username.
#
# Variables replaced by setup.sh:
#   __USER__ → current username
#   __WORKDIR__ → EvoMesh installation directory

[Unit]
Description=EvoMesh Multi-Role AI Orchestrator
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
ExecStart=/usr/bin/npx tsx --watch --watch-path=src bin/evomesh.ts serve --port 8123
Restart=always
RestartSec=5
User=__USER__
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
